# scripts/cloud-image/systemd/weknora.service
[Unit]
Description=WeKnora (docker compose)
Requires=docker.service
After=docker.service network-online.target weknora-firstboot.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/WeKnora
ExecStart=@DOCKER_BIN@ compose up -d
ExecStop=@DOCKER_BIN@ compose down
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
